# deploy/ssh-proxy.service
[Unit]
Description=SSH Proxy Core Server
Documentation=https://github.com/your-org/ssh-proxy-core
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ssh-proxy
Group=ssh-proxy
ExecStart=/usr/local/bin/ssh-proxy-core -c /etc/ssh-proxy/config.ini
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/log/ssh-proxy
ReadOnlyPaths=/etc/ssh-proxy

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
